# pi/crontab
@reboot /home/pi/i2c_server_starter.sh

